# portfolio.service
[Unit]
Description=Portfolio Website Next.js App
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/portfolio
Environment=NODE_ENV=production
Environment=PORT=3000
# Use standalone server for better performance
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=10
# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/portfolio

[Install]
WantedBy=multi-user.target
